# xdg-config/.config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comment=Simple X hotkey daemon
Exec=/usr/bin/sxhkd -c ~/.dotfiles/sxhkdrc
Terminal=false
Type=Application

